# material/profile.d/asdf.sh
: "${ASDF_DIR:="/etc/skel/.asdf"}"
: "${ASDF_DATA_DIR:="${WORKDIR}/.asdf"}"
: "${PATH:="/etc/skel/.asdf/bin"}"
export ASDF_DIR ASDF_DATA_DIR PATH
