# material/profile.d/asdf.sh
: "${ASDF_DIR:="/etc/skel/.asdf"}"
: "${ASDF_DATA_DIR:="${ASDF_DIR}"}"
export ASDF_DIR ASDF_DATA_DIR
