# material/profile.d/asdf.sh
: ${ASDF_DATA_DIR:="${WORKDIR}/.asdf"}

PATH="${ASDF_DATA_DIR}/shims:${PATH}"

