# material/profile.d/asdf.sh
: ${ASDF_DATA_DIR:="${WORKDIR}/.asdf"}

