# material/profile.d/config.sh
if [[ -d ~/.config ]]; then
    chown ${USER}: ~/.config
fi
