# material/profile.d/czsh.sh
# : "${CONTAINER_PROMPT:="placeholder"}"
# : "${CZSH_ISSUE_GH_REPO:="gh repo"}"
