# material/profile.d/go.sh
# : "${GOBIN:="$WORKDIR"}"
# : "${GOROOT:="/usr/local/go"}"
# # : "${GOPRIVATE:=""}"
# PATH=$PATH:${GOROOT}/bin:${GOBIN}
# export GOBIN GOROOT PATH GOPRIVATE

