# material/profile.d/nnn.sh
: "${SHELL:="czsh"}"
: "${EDITOR:="cvim"}"
: "${NNN_DE_FILE_MANAGER:="cvim"}"
: "${NNN_USE_EDITOR:="1"}"
export SHELL EDITOR NNN_DE_FILE_MANAGER NNN_USE_EDITOR
