# material/profile.d/nvm.sh
: "${NVM_DIR:="/etc/skel/.nvm"}"
export NVM_DIR
