# material/profile.d/oh-my-zsh.sh
[[ -e ~/.oh-my-zsh/cache ]] && sudo chown -R ${USER}: ~/.oh-my-zsh/cache
