# material/profile.d/oh-my-zsh.sh
sudo chown -R ${USER}. ~/.oh-my-zsh/cache
