# material/profile.d/pass.sh
: "${PASSWORD_STORE_DIR:="$WORKDIR"}"
export PASSWORD_STORE_DIR
