# material/profile.d/pass.sh
set -o allexport
: "${PASSWORD_STORE_DIR:="$WORKDIR"}"
set +o allexport
