# material/profile.d/python.sh
#: {PYTHONBREAKPOINT:="pudb.set_trace"}
: ${POETRY_VIRTUALENVS_IN_PROJECT:="true"}
