# material/profile.d/python.sh
PYENV_ROOT="${HOME}/.pyenv"
chmod a+w ${PYENV_ROOT} ${PYENV_ROOT}/shims ${PYENV_ROOT}/versions

