# material/profile.d/qt.sh
: "${QT_X11_NO_MITSHM:="1"}"
export QT_X11_NO_MITSHM
