# material/profile.d/qt.sh
#https://unix.stackexchange.com/questions/136606/mit-shared-memory-extension-or-mit-shm-over-ssh
: "${QT_X11_NO_MITSHM:="1"}"
