# material/profile.d/qt.sh
set -o allexport

: "${QT_X11_NO_MITSHM:="1"}"

set +o allexport
