# material/profile.d/snd.sh
if [[ ${USER} != "root" ]]; then
    usermod -aG audio ${USER}
fi
