# material/profile.d/terraform.sh
set -o allexport

: "${TF_INPUT:="0"}"

set +o allexport
