# material/profile.d/zsh.sh
LANG=C.UTF-8
LANGUAGE=C.UTF-8
LC_ALL=C.UTF-8
DISABLE_AUTO_UPDATE=true
SHELL=/usr/bin/zsh
export LANG LANGUAGE LC_ALL DISABLE_AUTO_UPDATE SHELL

chsh -s /bin/zsh $USER
