# material/profile.d/zsh.sh
set -o allexport

LANG=C.UTF-8
LANGUAGE=C.UTF-8
LC_ALL=C.UTF-8
DISABLE_AUTO_UPDATE=true
SHELL=/usr/bin/zsh

chsh -s /bin/zsh $USER

set +o allexport
