# material/profile.d/zsh.sh
LANG=C.UTF-8
LANGUAGE=C.UTF-8
LC_ALL=C.UTF-8
export LANG LANGUAGE LC_ALL

HISTFILE="${PWD}/.zsh_history"
export HISTFILE

chsh -s /bin/zsh $USER
